# recipes-core/auto-update/auto-update/71-automount.rules
# automount rules.
KERNEL=="sd[a-p][1-9]", ACTION=="add", RUN+="/bin/systemctl --no-block start device-automount@%k.service"
