# recipes-core/cangps-files/cangps-files/44-socketcan.rules
SUBSYSTEM=="net", KERNEL=="can[0-9]", RUN+="/bin/systemctl --no-block start socketcan@%k.service"
